# backend/uwsgi.service
[Unit]
Description=uWSGI instance to serve brighteyesbakend

[Service]
ExecStartPre=-/bin/bash -c 'mkdir -p /run/uwsgi; chown ubuntu:www-data /run/uwsgi'
ExecStart=/bin/bash -c 'cd /home/ubuntu/brighteyes/backend; source /home/ubuntu/env3/bin/activate; uwsgi --ini brighteyesbackend.ini --thunder-lock'

[Install]
WantedBy=multi-user.target
